# nix/pkgs/firecracker-rootfs/extra-files/etc/systemd/system/debug-shell.service
[Unit]
Description=Interactive debug shell
After=local-fs.target
Requires=local-fs.target

[Service]
Type=simple
ExecStart=/sbin/debug-shell
Restart=always
